# contrib/obs-mirror/systemd/obs-mirror.timer
[Unit]
Description=Run obs-mirror.service

[Timer]
OnCalendar=*-*-* 06:19:00

[Install]
WantedBy=timers.target
